# contrib/direwolf-dashboard.service
[Unit]
Description=Direwolf Dashboard - Live display of Direwolf activity
After=direwolf.service network.target
Wants=direwolf.service

[Service]
Type=simple
ExecStart=/home/pi/direwolf-dashboard/.venv/bin/direwolf-dashboard serve
Restart=always
RestartSec=5
User=pi
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/home/pi/direwolf-dashboard

[Install]
WantedBy=multi-user.target
